# work_shegnbao/ExpertContainer/supervisor/jisuandian-supervisord-production.conf
[unix_http_server]
file = /tmp/supervisor.sock
;chmod = 0777
;chown= root:felinx

;[inet_http_server]
# Web管理界面设定
;port=9001
;username = admin
;password = rnbtechrd

[supervisorctl]
; 必须和'unix_http_server'里面的设定匹配
serverurl = unix:///tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
user=root                 ; (default is current user, required if root)
;childlogdir=/var/log/supervisord/            ; ('AUTO' child log dir, default $TEMP)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; 管理的单个进程的配置，可以添加多个program


[program:alarm]
command=python3.4 /home/beop/ExpertContainer/AlarmSendNoticeQueue.py AlarmDataQueue_1 &
autorstart=true
autorestart=true


[program:expert1]
command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_1 &
autorstart=true
autorestart=true

[program:expert2]
command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_2 &
autorstart=true
autorestart=true

[program:expert3]
command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_3 &
autorstart=true
autorestart=true

[program:expert4]
command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_4 &
autorstart=true
autorestart=true

[program:expert5]
command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_5 &
autorstart=true
autorestart=true

[program:expert6]
command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_6 &
autorstart=true
autorestart=true

[program:expert7]
command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_7 &
autorstart=true
autorestart=true

[program:expert8]
command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_8 &
autorstart=true
autorestart=true

;[program:expert9]
;command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_9 &
;autorstart=true
;autorestart=true

;[program:expert10]
;command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_10 &
;autorstart=true
;autorestart=true

;[program:expert11]
;command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_11 &
;autorstart=true
;autorestart=true

;[program:expert12]
;command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_12 &
;autorstart=true
;autorestart=true

;[program:expert13]
;command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_13 &
;autorstart=true
;autorestart=true

;[program:expert14]
;command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_14 &
;autorstart=true
;autorestart=true

;[program:expert15]
;command=python3.4 /home/beop/ExpertContainer/RealCalcQueue.py RealCalcQueue_15 &
;autorstart=true
;autorestart=true

[program:diag1]
command=python3.4 /home/beop/ExpertContainer/DiagnosisCalcQueue.py DiagnosisCalcQueue_1 &
autorstart=true
autorestart=true

[program:diag2]
command=python3.4 /home/beop/ExpertContainer/DiagnosisCalcQueue.py DiagnosisCalcQueue_2 &
autorstart=true
autorestart=true

[program:diag3]
command=python3.4 /home/beop/ExpertContainer/DiagnosisCalcQueue.py DiagnosisCalcQueue_3 &
autorstart=true
autorestart=true

[program:diag4]
command=python3.4 /home/beop/ExpertContainer/DiagnosisCalcQueue.py DiagnosisCalcQueue_4 &
autorstart=true
autorestart=true

[program:force1]
command=python3.4 /home/beop/ExpertContainer/ForceCalcQueue.py ForceCalcQueue_1 &
autorstart=true
autorestart=true

[program:fixtime]
command=python3.4 /home/beop/ExpertContainer/FixedTimeCalc.py FixedTimeCalc.py &
autorstart=true
autorestart=true

[program:thirddata]
command=python3.4 /home/beop/ExpertContainer/updateThirdDataRealAndHistory.py updateThirdDataRealAndHistory.py &
autorstart=true
autorestart=true

[program:datataskproductor]
command=python3.4 /home/beop/ExpertContainer/DataTaskProductor.py DataTaskProductor.py &
autorstart=true
autorestart=true

[program:datataskconsumer]
command=python3.4 /home/beop/ExpertContainer/DataTaskConsumer.py DataTaskConsumer.py &
autorstart=true
autorestart=true

